# deploy/all-in-one/supervisord.conf
; service-console 单镜像进程托管:console(单一 FastAPI,hub+platform 已合并)+ nginx 两进程。
; 日志统一走容器 stdout(/dev/fd/1),`docker logs` 可见;任一进程崩溃自动重启。
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; console 先起(priority 小先启),nginx 随后反代。
[program:console]
command=/opt/run-console.sh
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=20
