# deploy/all-in-one/supervisord.conf
; service-console 单镜像进程托管:hub + platform + nginx 三进程。
; 日志统一走容器 stdout(/dev/fd/1),`docker logs` 可见;任一进程崩溃自动重启。
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; hub 先起(priority 小先启);platform 启动即调 hub,晚一点也无妨(调用侧容错)。
[program:hub]
command=/opt/run-hub.sh
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=10

[program:platform]
command=/opt/run-platform.sh
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=30
